# acdc.service
[Unit]
Description=ACDC Gunicorn Application Server
After=network.target postgresql.service redis.service

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/nusu85_new
EnvironmentFile=/home/ubuntu/nusu85_new/.env
ExecStart=/home/ubuntu/venv/bin/gunicorn \
    acdc_config.wsgi:application \
    --workers 4 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --max-requests 1000 \
    --access-logfile /var/log/acdc/access.log \
    --error-logfile /var/log/acdc/error.log
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
